# garuda-ai-optimizations/usb-performance-optimization.service
[Unit]
Description=USB Performance Optimization for Rescue System
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-rescue-optimizer.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
